# config/personal-site.conf
description "Gunicorn application server running personal site"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data

chdir /home/ubuntu/personal-site
exec venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8003 personal-site.wsgi:application
